# files/files_etcdRAM/usr/lib/systemd/system/kubelet.service.d/etcdRAM-restore-etcd-pre-start.conf

[Unit]
After=containerd.service

[Service]
ExecStartPre=-/opt/etcdRAM/bin/etcd-restore-to-tmpfs
