# files/files_etcdRAM/usr/lib/systemd/system/kubelet.service.d/etcdRAM-backup-etcd-post-stop.conf

[Service]
ExecStop=/opt/etcdRAM/bin/etcd-cron
