# files/files_etcdRAM/etc/cron.d/etcdRAM-etcd
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
HOME=/root
MAILTO=root

*/20 * * * * root /opt/etcdRAM/bin/etcd-cron > /var/log/etcdRAM-etcd-backup 2>&1
30 0 * * * root /opt/etcdRAM/bin/etcdctl defrag --cluster > /var/log/etcdRAM-etcd-defrag 2>&1
